# dashcam/package/gnss-logger/gnss-logger.service
[Unit]
Description=GNSS Logger
After=network.target

[Service]
WorkingDirectory=/opt/dashcam/bin
ExecStartPre=/bin/mkdir -p /mnt/data/gps
ExecStart=/opt/dashcam/bin/gnsslogger -mga-offline-file-path /opt/dashcam/bin/mgaoffline.ubx
Restart=always

[Install]
WantedBy=multi-user.target
