# dashcam/package/camera-bridge/camera-bridge.timer
[Unit]
Description=Timer to start camera-bridge

[Timer]
OnStartupSec=60sec

[Install]
WantedBy=timers.target
