# dashcam/package/camera-bridge/camera-bridge.timer
[Unit]
Description=Timer to start camera-bridge

[Timer]
OnStartupSec=5min

[Install]
WantedBy=timers.target