# dashcam/board/raspberrypi/overlays/common/etc/systemd/system/bootbit.service
[Unit]
Description=Dashcam Boot BIT
After=rauc.service
Requires=rauc.service

[Service]
Type=simple
ExecStart=/usr/bin/bootbit
User=root
