# dashcam/board/raspberrypi/overlays/common/etc/systemd/system/wifiP2P.service
[Unit]
Description=WiFi P2P auto connector
After=network.target dnsmasq.service dhcpcd.service camera-node.service

[Service]
Type=simple
PIDFile=/run/wifiP2P.pid
ExecStart=/opt/dashcam/bin/network/wifiP2Pdaemon.sh
Restart=always

[Install]
WantedBy=multi-user.target
